# packaging/rpm/systemd/sauron-ingest.service
[Unit]
Description=Sauron ingest edge and worker pool
Documentation=file:///usr/share/doc/sauron-server/SETUP.md
After=network-online.target sauron-migrate.service
Wants=network-online.target

[Service]
Type=exec
User=sauron
Group=sauron
EnvironmentFile=/etc/sauron/sauron.env
EnvironmentFile=-/etc/sauron/ingest.env
ExecStart=/usr/bin/sauron-ingest
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
StateDirectory=sauron
StateDirectoryMode=0750
ReadWritePaths=/var/lib/sauron

[Install]
WantedBy=multi-user.target
